# gaimoded/gaimoded.service
[Unit]
Description=Gaimode service for optimizing process performance

[Service]
ExecStart=/usr/bin/gaimoded

[Install]
WantedBy=multi-user.target
